# systemd/nstr-report-repost.timer
[Unit]
Description=Repost nstr-report summary at 9am and 12pm ET
Documentation=https://github.com/EthnTuttle/nstr-report

[Timer]
# Repost at 9am and 12pm Eastern Time
OnCalendar=*-*-* 09:00:00 America/New_York
OnCalendar=*-*-* 12:00:00 America/New_York
RandomizedDelaySec=120
Persistent=true

[Install]
WantedBy=timers.target
